# loghome-agent/supervisor_config.ini
[program:loghome_simulator]
command=python simulator.py --forever
directory=/d:/LogHome/loghome-agent
autostart=true
autorestart=true
startsecs=5
startretries=3
redirect_stderr=true
stdout_logfile=/d:/LogHome/loghome-agent/logs/simulator.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/d:/LogHome/loghome-agent/logs/simulator_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5 